# picec.service
[Unit]
Description=Cursor control via CEC

[Service]
Type=simple
ExecStart=/bin/bash -lc picec
ExecReload=/bin/kill -USR1 $MAINPID

[Install]
WantedBy=default.target
