# picec.service
[Unit]
Description=Cursor control via CEC

[Service]
Type=simple
ExecStart=/bin/bash -lc picec.py

[Install]
WantedBy=default.target
